# openstack/python-heat/openstack-heat/centos/files/openstack-heat-api-cloudwatch.service
[Unit]
Description=OpenStack Heat CloudWatch API Service
After=syslog.target network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/heat-api-cloudwatch --config-file /usr/share/heat/heat-dist.conf --config-file /etc/heat/heat.conf

[Install]
WantedBy=multi-user.target
