# openstack/python-nova/centos/files/openstack-nova-metadata-api.service
[Unit]
Description=OpenStack Nova Metadata API Server
After=syslog.target network.target

[Service]
Type=simple
#NotifyAccess=all
TimeoutStartSec=0
#Restart=always
User=root
ExecStart=/usr/bin/nova-api-metadata

[Install]
WantedBy=multi-user.target

