# openstack/openstack-magnum/centos/files/openstack-magnum-api.service
[Unit]
Description=OpenStack Magnum API Service
After=syslog.target network.target

[Service]
Type=simple
User=magnum
ExecStart=/usr/bin/magnum-api
PrivateTmp=true
NotifyAccess=all
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
